# config/etc/supervisor/conf.d-temp/php-fpm.conf
# PHP FPM
[program:php-fpm]
command=php-fpm8.4 -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startretries=0